# digital-intern/digital-intern-public-dashboard.service
[Unit]
Description=Digital Intern Public Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zeph
WorkingDirectory=/home/zeph/trading-intelligence/digital-intern
EnvironmentFile=-/home/zeph/trading-intelligence/digital-intern/.env
Environment="PATH=/home/zeph/.local/bin:/home/zeph/.nvm/versions/node/v24.15.0/bin:/usr/local/bin:/usr/bin:/bin"
Environment="WEB_SERVER_HOST=0.0.0.0"
Environment="WEB_SERVER_PORT=8080"
ExecStart=/usr/bin/python3 -m dashboard.web_server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
